# nginx-yum-repo/rebuild-repo.path
[Unit]
Description= Rebuild repo

[Path]
PathModified=/srv/data/repoflags

[Install]
WantedBy=multi-user.target
