# nginx-yum-repo/rebuild-repo.service
[Unit]
Description= Rebuild repo

[Service]
Type=oneshot
User=nginx
Group=nginx
ExecStart=/opt/bin/rebuild_repo.sh
